# scripts/voice/launchd/com.livhana.voice.tts.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.livhana.voice.tts</string>
  
  <key>ProgramArguments</key>
  <array>
    <string>/usr/bin/env</string>
    <string>bash</string>
    <string>-lc</string>
    <string>/Users/jesseniesen/LivHana-Trinity-Local/LivHana-SoT/scripts/voice/start_kokoro_tts.sh</string>
  </array>
  
  <key>WorkingDirectory</key>
  <string>/Users/jesseniesen/LivHana-Trinity-Local/LivHana-SoT</string>
  
  <key>RunAtLoad</key>
  <true/>
  
  <key>KeepAlive</key>
  <dict>
    <key>SuccessfulExit</key>
    <false/>
  </dict>
  
  <key>EnvironmentVariables</key>
  <dict>
    <key>PORT</key>
    <string>8880</string>
    <key>HOST</key>
    <string>127.0.0.1</string>
  </dict>
  
  <key>StandardOutPath</key>
  <string>/Users/jesseniesen/Library/Logs/LivHana/voice-tts.out.log</string>
  
  <key>StandardErrorPath</key>
  <string>/Users/jesseniesen/Library/Logs/LivHana/voice-tts.err.log</string>
</dict>
</plist>

